# unused/jehon-system-kiosk/lib/udev/rules.d/98-kiosk-smartcard-change.rules

SUBSYSTEM=="block", ACTION=="change", ENV{DISK_MEDIA_CHANGE}=="1" RUN+="/usr/bin/kiosk-smartcard-change.sh"
